# collectd/collectd-wrapper.service
[Unit]
Description=collectd-wrapper daemon

[Service]
ExecStart=/root/collectd/script.sh
StandardOutput=null

[Install]
WantedBy=multi-user.target
